# deploy/ntfy-voipms.service
[Unit]
Description=ntfy push server (voipms-sms) — private, LAN/VPN
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# User-space ntfy on port 8090, file-based cache (no /var, no root).
ExecStart=%h/.local/bin/ntfy serve --listen-http :8098 --cache-file %h/.local/share/ntfy/cache.db --base-url http://192.168.1.12:8098
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
